# 服务集群安装/supervisor安装/nta-web-base.ini
[program:nta-web-base]
command=/opt/nta/web/base/bin/startup
directory=/opt/nta/web/base
priority=60
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/data/logs/web/base/base.stdout.log
stderr_logfile=/data/logs/web/base/base.stderr.log
